# resources/keywork-xdg-autostart.target
[Unit]
Description=XDG autostart for the Keywork session
BindsTo=xdg-desktop-autostart.target
PartOf=keywork-session.target
After=keywork-session.target
